# root/etc/services.d/resilio-sync/run
#!/usr/bin/with-contenv sh

APP_PATH="/app/resilio-sync" 
CONF_PATH="/config/resilio-sync.conf"

umask 022

exec \
	s6-setuidgid abc "$APP_PATH"/resilio --config "$CONF_PATH" --nodaemon